# crontab
0 3 * * * root python /run_letsencrypt.py
